# docs/systemd_du-wechat-ilink.service
[Unit]
Description=Du Gateway - WeChat iLink Connector
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/du-gateway/connectors/wechat_ilink
EnvironmentFile=/root/du-gateway/connectors/wechat_ilink/.env
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=2

# 连接器会写本地状态文件（含 bot_token / get_updates_buf），建议限制权限
UMask=0077

[Install]
WantedBy=multi-user.target

